# ciao/stock/deploy/com.ciao.server.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ciao.server</string>

    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>local.ciaobot.app</string>
    </array>

    <key>ProgramArguments</key>
    <array>
        <string>{{CIAO_PYTHON}}</string>
        <string>-m</string>
        <string>ciao.cli</string>
        <string>run</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{CIAO_WORKSPACE}}</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CIAO_WORKSPACE</key>
        <string>{{CIAO_WORKSPACE}}</string>
        <key>CIAO_PORT</key>
        <string>{{CIAO_PORT}}</string>
        <key>PATH</key>
        <string>{{CIAO_PATH}}</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{{CIAO_WORKSPACE}}/.runtime/ciao.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{CIAO_WORKSPACE}}/.runtime/ciao.stderr.log</string>
</dict>
</plist>
